# vms/fw01/iptables/adamo.service
[Unit]
Description=ADAMO

[Service]
ExecStart=/usr/local/bin/adamo.sh

[Install]
WantedBy=multi-user.target
